# scripts/gpu-shutdown.service
[Unit]
Description=GPU PC Shutdown Server
After=network.target

[Service]
Type=simple
User=arnab
ExecStart=/usr/bin/python3 /home/arnab/shutdown-server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
